# sw/OPAE/driver/opae-intel-fpga-driver-1.3.0-2/40-intel-fpga.rules
KERNEL=="intel-fpga*[0-9]*", ACTION="add", GROUP="root", MODE="0664", RUN+="/bin/find /sys/class/fpga/intel-fpga-dev*/ -type d -execdir /bin/chown root:root {} \;", RUN+="/bin/find /sys/class/fpga/intel-fpga-dev*/ -type d -execdir /bin/chmod 755 {} \;", RUN+="/bin/find /sys/class/fpga/intel-fpga-dev*/ -iname freeze -execdir /bin/chown root:root {} \;", RUN+="/bin/find /sys/class/fpga/intel-fpga-dev*/ -iname freeze -execdir /bin/chmod 666 {} \;"
KERNEL=="intel-fpga*[0-9]*", ACTION="change", GROUP="root", ENV{TRIGGER}!="none", MODE="0664", RUN+="/bin/find /sys/class/fpga/intel-fpga-dev*/ -type d -execdir /bin/chown root:root {} \;", RUN+="/bin/find /sys/class/fpga/intel-fpga-dev*/ -type d -execdir /bin/chmod 755 {} \;", RUN+="/bin/find /sys/class/fpga/intel-fpga-dev*/ -iname freeze -execdir /bin/chown root:root {} \;", RUN+="/bin/find /sys/class/fpga/intel-fpga-dev*/ -iname freeze -execdir /bin/chmod 666 {} \;"
